# debian/amavisd-new.amavisd-snmp-subagent.service
[Unit]
Description=Exports amavis SNMP data
Documentation=http://www.ijs.si/software/amavisd/#doc
After=network.target amavis.service

[Service]
User=amavis
Group=amavis
ExecStart=/usr/sbin/amavisd-snmp-subagent -f
Restart=on-failure
PrivateTmp=true
CapabilityBoundingSet=
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
